# templates/phd_rhel.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          phd
# Required-Start:    $remote_fs $syslog mysql
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Wrapper around phabricator's phd daemon launcher
### END INIT INFO

DESC="phd daemon launcher"
NAME=phd
DAEMON=<%= @phd_path %>/$NAME
SCRIPTNAME=/etc/init.d/$NAME

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

RETVAL=0
case "$1" in
    start)
        [ "$VERBOSE" != no ] && echo "Starting $DESC" "$NAME"
        "$DAEMON" start
        RETVAL="$?"
        case "$?" in
            0) [ "$VERBOSE" != no ] && echo 0 ;;
            1) [ "$VERBOSE" != no ] && echo 255 ;; # Already started
            *) [ "$VERBOSE" != no ] && echo 1 ;; # Failed to start
        esac
        ;;
    stop)
        [ "$VERBOSE" != no ] && echo "Stopping $DESC" "$NAME"
        "$DAEMON" stop
        RETVAL="$?"
        case "$?" in
            0) [ "$VERBOSE" != no ] && echo 0 ;;
            *) [ "$VERBOSE" != no ] && echo 1 ;; # Failed to stop
        esac
        ;;
    status)
	"$DAEMON" status && exit 0 || exit $?
	;;
    restart)
        echo "Restarting $DESC" "$NAME"
        "$DAEMON" restart
        RETVAL="$?"
        case "$?" in
            0) echo 0 ;;
            *) echo 1 ;; # Failed to restart
        esac
        ;;
    *)
        echo "Usage: $SCRIPTNAME {start|stop|status|restart}" >&2
        exit 3
        ;;
esac

exit "$RETVAL"